# deploy/systemd/egi-web-monitoring-db-backup.service
[Unit]
Description=Backup EGI Web Monitoring PostgreSQL database
After=postgresql@16-main.service
Requires=postgresql@16-main.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/egi-web-monitoring-db-backup
ProtectSystem=full
PrivateTmp=true
